# packaging/rookery.container
# Rookery, dogfooded as a Quadlet. Install to /etc/containers/systemd/
# (rootful) or ~/.config/containers/systemd/ (rootless), then:
#   systemctl daemon-reload && systemctl start rookery
#
# Note: running in a container limits Rookery to the mounted unit dirs and
# requires the journal and systemd sockets below. The simpler install is the
# bare binary with the sibling rookery.service unit.

[Unit]
Description=Rookery — Quadlet-native web UI for Podman
Wants=network-online.target
After=network-online.target

[Container]
Image=ghcr.io/tobagin/rookery:latest
# Bind all interfaces inside the container; the host side stays loopback-only.
Environment=ROOKERY_LISTEN=0.0.0.0:7665
PublishPort=127.0.0.1:7665:7665
Volume=/etc/containers/systemd:/etc/containers/systemd
Volume=/run/systemd:/run/systemd
# systemctl show/list talk to systemd over the D-Bus system bus, not the
# private socket. Mount the directory, not the socket file, so the mount
# survives a dbus restart. Not :ro — connecting to a socket needs write.
Volume=/run/dbus:/run/dbus
Volume=/var/log/journal:/var/log/journal:ro
Volume=/run/podman/podman.sock:/run/podman/podman.sock:ro
# The HOST's quadlet generator, so validation matches the Podman that will
# actually run the units (Fedora/RHEL path — adjust for other distros).
# Without it, validation degrades to "unavailable" but everything else works.
Volume=/usr/libexec/podman/quadlet:/usr/libexec/podman/quadlet:ro
SecurityLabelDisable=true

[Service]
Restart=always

[Install]
WantedBy=multi-user.target default.target
